# supervisord.conf
[supervisord]
nodaemon=true

[program:vm]
command=/start.sh
autorestart=true
