# supervisord.conf
[supervisord]
nodaemon=true

[program:android]
command=/start.sh
autorestart=true
